# church_is.conf
[program:church_is]
command=/home/nanoafrika/church_is_supervisor
user=nanoafrika
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/home/nanoafrika/logs/church_is.log